# system_files/usr/lib/systemd/user/swaync.service.d/set-config.conf
[Service]
ExecStart=
ExecStart=/usr/bin/swaync -s $SWAYNC_CONFIG_STYLE -c $SWAYNC_CONFIG_FILE